# todoapp.socket
[Socket]
ListenStream=80
Service=todoapp.service