# todoapp.socket
[Unit]
Description=My To-Do app socket

[Socket]
ListenStream=80